# UbuntuDev/createuser.service
[Unit]
Description="Create the user"

[Service]
type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "sbin/createuser.sh"
